# service/logcollector.service.template
[Unit]
Description=Log Collector Daemon
After=network.target

[Service]
User={USER}
ExecStart={EXEC_START}
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
